# deploy/systemd/devops-hub.service
[Unit]
Description=DevOps Hub
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/srv/devops-hub
EnvironmentFile=/srv/devops-hub/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
